# ops/breakout_1425.cron.example
# 14:25 突破提醒任务（北京时间，周一到周五）
# 建议先确认 .env 中 BREAKOUT_WECHAT_WEBHOOK_URL / ALERT_WEBHOOK_URL 已配置
25 14 * * 1-5 cd /Users/alvachen/aiproject/tradingart && /bin/bash run_breakout_1425.sh
